# email_parser/crontab
# START CRON JOB
*/5 * * * * python3 /app/retrieve_email.py > /var/log/script.log
# END CRON JOB